# supervisord.conf
[program:node]
command = node server.js
directory = /home/dotcloud/current
environment = NODE_ENV=dotcloud
